# .service
[Unit]
Description=Gunicorn instance to serve GumGum-Flask-Application
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/opt/releases/GumGum-Flask-Application
Environment="PATH=/opt/releases/GumGum-Flask-Application/env/bin"
ExecStart=/opt/releases/GumGum-Flask-Application/env/bin/gunicorn --workers 3 --bind unix:GumGum-Flask-Application.sock -m 007 wsgi:app

[Install]
WantedBy=multi-user.target
